# stow/launchd/Library/LaunchAgents/local.borgmatic.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>local.borgmatic</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/borgmatic</string>
        <string>--syslog-verbosity</string>
        <string>1</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/mike</string>
    </dict>
    <key>StartInterval</key>
    <integer>86400</integer>
</dict>
</plist>
