# stow/launchd/Library/LaunchAgents/org.conigliaro.ghm.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>org.conigliaro.ghm</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/ghm</string>
        <string>/Users/mike/Documents/Backups/GitHub</string>
        <string>--token</string>
        <string>file:///Users/mike/.config/github/token</string>
        <string>--exclude-forks</string>
        <string>--exclude-owner</string>
        <string>(ganglia|monitoringsucks)</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>StandardOutPath</key>
    <string>/tmp/org.conigliaro.ghm.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/org.conigliaro.ghm.log</string>
</dict>
</plist>
